# packaging/systemd/plocate-server-updatedb.service
[Unit]
Description=Refresh the plocate-server index (via the API)
Wants=plocate-server.service
After=plocate-server.service

[Service]
Type=oneshot

# Ask the running plocate-server to rebuild its index. The server is the sole
# executor of updatedb (it owns the reindex lock, the stat cache invalidation,
# and the last_run record), so routing the timer through the HTTP endpoint
# guarantees the timer and the maintenance UI button can never race.
ExecStart=/usr/bin/curl -fsS -X POST --max-time 30 http://127.0.0.1:8787/api/reindex

# Best-effort: if the server is down, the curl fails and the timer retries on
# its next schedule. Nothing to harden here — no filesystem access, no caps.

PrivateTmp=true
