# packaging/systemd/plocate-server.service
[Unit]
Description=plocate-server (plocate-backed filename search REST API)
Documentation=https://github.com/fanyang89/plocate-server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=plocate-server
Group=plocate-server

# All configuration lives in /etc/plocate-server/plocate-server.env (the single
# source — no inline Environment= defaults to drift out of sync).
EnvironmentFile=-/etc/plocate-server/plocate-server.env

ExecStart=/usr/bin/plocate-server

Restart=on-failure
RestartSec=3

# updatedb must traverse the whole tree to index it; this capability lets the
# unprivileged service read past file-mode bits without running as root.
AmbientCapabilities=CAP_DAC_READ_SEARCH
CapabilityBoundingSet=CAP_DAC_READ_SEARCH

# The server keeps no index in RAM; 1G is ample headroom.
MemoryMax=1G
MemoryHigh=800M

# Yield to the foreground service under contention.
CPUWeight=20
Nice=19
IOSchedulingClass=idle
IOSchedulingPriority=7

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
ReadWritePaths=/var/lib/plocate-server

[Install]
WantedBy=multi-user.target
