# packaging/systemd/plocate-server-updatedb.service
[Unit]
Description=Refresh the plocate-server index
After=plocate-server.service
Wants=plocate-server.service

[Service]
Type=oneshot
User=plocate-server
Group=plocate-server

Environment=PLOCATE_SERVER_BASE_PATH=/srv/files
Environment=PLOCATE_SERVER_DB_PATH=/var/lib/plocate-server/files.db
EnvironmentFile=-/etc/plocate-server/plocate-server.env

# Rebuild the on-disk trigram index. ${...} expands from Environment/EnvironmentFile.
ExecStart=/usr/libexec/plocate-server/updatedb \
    -U ${PLOCATE_SERVER_BASE_PATH} \
    -o /var/lib/plocate-server/files.db \
    --require-visibility no

# Read the whole tree without root, at lowest priority so the foreground
# service is never disturbed.
AmbientCapabilities=CAP_DAC_READ_SEARCH
CapabilityBoundingSet=CAP_DAC_READ_SEARCH
Nice=19
IOSchedulingClass=idle
IOSchedulingPriority=7

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateNetwork=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
ReadWritePaths=/var/lib/plocate-server
