# deploy/plocate-server.service
[Unit]
Description=plocate-server (plocate-backed filename search REST API)
After=network-online.target
Wants=network-online.target
# Requires plocate (provides plocate + updatedb) installed on the host.

[Service]
Type=simple
User=plocate-server
Group=plocate-server

# --- Override ExecStart via `systemctl edit plocate-server` for your tree ---
# Periodic index refresh is NOT configured in this bare-deploy example.
# Either set up a systemd timer / cron that POSTs to /api/reindex, or trigger
# refresh manually via the maintenance dialog in the web UI.
ExecStart=/usr/local/bin/plocate-server \
    --base-path=/srv/files \
    --bind=127.0.0.1:8787 \
    --db-path=/var/lib/plocate-server/files.db
# -----------------------------------------------------------------------

Restart=on-failure
RestartSec=3

# ===== Permissions: updatedb must traverse the whole tree to index it =====
# CAP_DAC_READ_SEARCH lets the service (and its updatedb child) read past
# file-mode bits, so the index is complete without running as root.
AmbientCapabilities=CAP_DAC_READ_SEARCH
CapabilityBoundingSet=CAP_DAC_READ_SEARCH

# ===== Memory: 1G is ample; the server keeps no index in RAM (plocate's
#       on-disk trigram index is queried by a short-lived child process) =====
MemoryMax=1G
MemoryHigh=800M

# ===== CPU: yield under contention, capped at 2 cores =====
CPUWeight=20
# Hard ceiling — reindex never uses more than 2 cores, even on an idle host.
CPUQuota=200%

# ===== Scheduler priority + disk IO: always defer to the foreground service =====
Nice=19
IOSchedulingClass=idle
IOSchedulingPriority=7

# ===== Basic hardening =====
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
# The db-dir must be writable (updatedb writes the index there). The indexed
# base-path is read-only and read via the capability above.
ReadWritePaths=/var/lib/plocate-server

[Install]
WantedBy=multi-user.target
